# iso/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
# Override getty to autologin the minecrarch user without password.
# This is intentional for an appliance-style system (ADR-0006: no desktop,
# no login screen — boot goes directly to Gamescope session).
ExecStart=
ExecStart=-/sbin/agetty --autologin minecrarch --noclear %I $TERM
